# lib/upstart/foreman-APP-N.conf
#{{{application}}}-{{{process}}}-{{{number}}}

start on starting {{{application}}}-{{{process}}}
stop  on stopping {{{application}}}-{{{process}}}
respawn

exec su - {{{user}}} -c 'cd {{{cwd}}}; {{#envs}}export {{{key}}}={{{value}}}; {{/envs}}{{{command}}} >> {{{logs}}}/{{{application}}}/{{{process}}}-{{{number}}}.log 2>&1'

post-stop script
    mv {{{logs}}}/{{{application}}}/{{{process}}}-{{{number}}}.log {{{logs}}}/{{{application}}}/{{{process}}}-{{{number}}}.`date -u +%Y-%m-%dT%T.%3NZ`.log
end script